# llama-server.service
[Unit]
Description=Llama Server (RAYD AI)
After=network.target

[Service]
Type=simple
User=rami
ExecStart=/opt/bitnet/build/bin/llama-server \
  -m /opt/bitnet/models/BitNet-b1.58-2B-4T/ggml-model-i2_s.gguf \
  -t 14 -tb 14 --host 0.0.0.0 --port 8081 -c 4096 --parallel 2
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
