# llama-server.service
[Unit]
Description=Llama Server (RAYD AI)
After=network.target

[Service]
Type=simple
User=stats
ExecStart=/home/stats/BitNet/build/bin/llama-server \
  -m /home/stats/BitNet/models/Meta-Llama-3.1-8B-Instruct-Q4_K_M.gguf \
  -t 6 --host 0.0.0.0 --port 8081 -c 2048
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
